# 18-c0-4d-e0-80-28/etc/systemd/system/var-lib-containers.mount
[Unit]
Requires=systemd-fsck@dev-md-containers.service
After=systemd-fsck@dev-md-containers.service

[Mount]
Where=/var/lib/containers
What=/dev/md/containers
Type=xfs

[Install]
RequiredBy=local-fs.target