# global/templates/reboot.service
[Unit]
Description=Reboot server

[Service]
Type=oneshot
ExecStart=/usr/sbin/reboot