# 18-c0-4d-e0-80-28/templates/systemd/nslookup.service
[Unit]
Description=Wait for DNS resolution
After=nss-lookup.target
Requires=nss-lookup.target

[Service]
Type=oneshot
ExecStart=/usr/bin/nslookup {{ .values.hostname }}.studiop

[Install]
RequiredBy=docker.service